# deploy/diarized-transcriber.container
[Unit]
Description=Diarized Transcriber Service
After=redis.service

[Container]
Image=localhost/diarized-transcriber:latest
ContainerName=diarized-transcriber
PublishPort=127.0.0.1:8721:8000
Network=gpu-services
AddDevice=nvidia.com/gpu=all
Environment=NVIDIA_VISIBLE_DEVICES=all

Environment=DIARIZED_REDIS_URL=redis://redis:6379
Environment=DIARIZED_MODEL_SIZE=large-v3-turbo
Environment=DIARIZED_PORT=8000
Volume=%h/.cache/huggingface:/root/.cache/huggingface:Z
Volume=%h/.config/huggingface:/root/.config/huggingface:ro,Z
Volume=diarized-speaker-db:/data/speaker_profiles_db:Z
Volume=/tmp/diarized-transcriber:/tmp:Z

Environment=DIARIZED_SPEAKER_DB_PATH=/data/speaker_profiles_db

[Service]
Restart=always
TimeoutStartSec=180
TimeoutStopSec=60

[Install]
WantedBy=default.target
